# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:backend]
directory=/opt/gemini-h5/backend
command=/usr/bin/python3 main.py
autorestart=true
stdout_logfile=/var/log/backend.log
stderr_logfile=/var/log/backend.err.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
stdout_logfile=/var/log/nginx-supervisor.log
stderr_logfile=/var/log/nginx-supervisor.err.log

[program:frpc]
command=/opt/frp/frpc -c /opt/frp/frpc.ini
autorestart=true
stdout_logfile=/var/log/frpc.log
stderr_logfile=/var/log/frpc.err.log
